# dev/updated_rules/startup_folder_file_write.json
[
    {
        "group": "Sigma",
        "kind": "individual",
        "document": {
            "kind": "jsonl",
            "data": {
                "@timestamp": "2022-09-27T12:23:35.784Z",
                "agent": {
                    "ephemeral_id": "f3e67ec5-847a-448b-8ec4-bac040a44990",
                    "hostname": "CLIENT1",
                    "id": "cd699c72-0f46-4cd9-938a-43faf035a44a",
                    "name": "CLIENT1",
                    "type": "winlogbeat",
                    "version": "7.10.2"
                },
                "ecs": {
                    "version": "1.5.0"
                },
                "event": {
                    "action": "File created (rule: FileCreate)",
                    "category": [
                        "file"
                    ],
                    "code": 11,
                    "created": "2022-09-27T12:27:28.040Z",
                    "kind": "event",
                    "module": "sysmon",
                    "provider": "Microsoft-Windows-Sysmon",
                    "type": [
                        "creation"
                    ]
                },
                "file": {
                    "directory": "C:\\Users\\setup\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs",
                    "name": "Startup",
                    "path": "C:\\Users\\setup\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup"
                },
                "host": {
                    "name": "CLIENT"
                },
                "log": {
                    "level": "information"
                },
                "message": "File created:\nRuleName: T1023\nUtcTime: 2022-09-27 12:23:35.784\nProcessGuid: {0e37ca90-eb41-6332-4400-000000000300}\nProcessId: 3792\nImage: C:\\Windows\\Explorer.EXE\nTargetFilename: C:\\Users\\setup\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\nCreationUtcTime: 2022-09-27 12:23:35.784\nUser: CLIENT\\setup",
                "process": {
                    "entity_id": "{0e37ca90-eb41-6332-4400-000000000300}",
                    "executable": "C:\\Windows\\Explorer.EXE",
                    "name": "Explorer.EXE",
                    "pid": 3792
                },
                "rule": {
                    "name": "T1023"
                },
                "winlog": {
                    "api": "wineventlog",
                    "channel": "Microsoft-Windows-Sysmon/Operational",
                    "computer_name": "CLIENT",
                    "event_data": {
                        "CreationUtcTime": "2022-09-27 12:23:35.784",
                        "User": "CLIENT\\setup"
                    },
                    "event_id": 11,
                    "opcode": "Info",
                    "process": {
                        "pid": 1328,
                        "thread": {
                            "id": 2576
                        }
                    },
                    "provider_guid": "{5770385f-c22a-43e0-bf4c-06f5698ffbd9}",
                    "provider_name": "Microsoft-Windows-Sysmon",
                    "record_id": 686,
                    "task": "File created (rule: FileCreate)",
                    "user": {
                        "domain": "NT AUTHORITY",
                        "identifier": "S-1-5-18",
                        "name": "SYSTEM",
                        "type": "User"
                    },
                    "version": 2
                }
            }
        },
        "name": "Startup Folder File Write",
        "timestamp": "2022-09-27T12:23:35.784+00:00",
        "authors": [
            "Roberto Rodriguez (Cyb3rWard0g)",
            "OTR (Open Threat Research)"
        ],
        "level": "medium",
        "source": "sigma",
        "status": "experimental",
        "falsepositives": [
            "An FP could be caused by legitimate application writing shortcuts for example. This folder should always be inspected to make sure that all the files in there are legitimate"
        ],
        "id": "2aa0a6b4-a865-495b-ab51-c28249537b75",
        "logsource": {
            "category": "file_event",
            "product": "windows"
        },
        "references": [
            "https://github.com/OTRF/detection-hackathon-apt29/issues/12",
            "https://threathunterplaybook.com/evals/apt29/detections/5.B.1_611FCA99-97D0-4873-9E51-1C1BA2DBB40D.html"
        ],
        "tags": [
            "attack.persistence",
            "attack.t1547.001"
        ]
    },
    {
        "group": "Sigma",
        "kind": "individual",
        "document": {
            "kind": "jsonl",
            "data": {
                "@timestamp": "2022-09-27T12:23:35.784Z",
                "agent": {
                    "ephemeral_id": "f3e67ec5-847a-448b-8ec4-bac040a44990",
                    "hostname": "CLIENT1",
                    "id": "cd699c72-0f46-4cd9-938a-43faf035a44a",
                    "name": "CLIENT1",
                    "type": "winlogbeat",
                    "version": "7.10.2"
                },
                "ecs": {
                    "version": "1.5.0"
                },
                "event": {
                    "action": "File created (rule: FileCreate)",
                    "category": [
                        "file"
                    ],
                    "code": 11,
                    "created": "2022-09-27T12:27:28.040Z",
                    "kind": "event",
                    "module": "sysmon",
                    "provider": "Microsoft-Windows-Sysmon",
                    "type": [
                        "creation"
                    ]
                },
                "file": {
                    "directory": "C:\\Users\\setup\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup",
                    "extension": "ini",
                    "name": "desktop.ini",
                    "path": "C:\\Users\\setup\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\desktop.ini"
                },
                "host": {
                    "name": "CLIENT"
                },
                "log": {
                    "level": "information"
                },
                "message": "File created:\nRuleName: T1023\nUtcTime: 2022-09-27 12:23:35.784\nProcessGuid: {0e37ca90-eb41-6332-4400-000000000300}\nProcessId: 3792\nImage: C:\\Windows\\Explorer.EXE\nTargetFilename: C:\\Users\\setup\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\desktop.ini\nCreationUtcTime: 2022-09-27 12:23:35.784\nUser: CLIENT\\setup",
                "process": {
                    "entity_id": "{0e37ca90-eb41-6332-4400-000000000300}",
                    "executable": "C:\\Windows\\Explorer.EXE",
                    "name": "Explorer.EXE",
                    "pid": 3792
                },
                "rule": {
                    "name": "T1023"
                },
                "winlog": {
                    "api": "wineventlog",
                    "channel": "Microsoft-Windows-Sysmon/Operational",
                    "computer_name": "CLIENT",
                    "event_data": {
                        "CreationUtcTime": "2022-09-27 12:23:35.784",
                        "User": "CLIENT\\setup"
                    },
                    "event_id": 11,
                    "opcode": "Info",
                    "process": {
                        "pid": 1328,
                        "thread": {
                            "id": 2576
                        }
                    },
                    "provider_guid": "{5770385f-c22a-43e0-bf4c-06f5698ffbd9}",
                    "provider_name": "Microsoft-Windows-Sysmon",
                    "record_id": 687,
                    "task": "File created (rule: FileCreate)",
                    "user": {
                        "domain": "NT AUTHORITY",
                        "identifier": "S-1-5-18",
                        "name": "SYSTEM",
                        "type": "User"
                    },
                    "version": 2
                }
            }
        },
        "name": "Startup Folder File Write",
        "timestamp": "2022-09-27T12:23:35.784+00:00",
        "authors": [
            "Roberto Rodriguez (Cyb3rWard0g)",
            "OTR (Open Threat Research)"
        ],
        "level": "medium",
        "source": "sigma",
        "status": "experimental",
        "falsepositives": [
            "An FP could be caused by legitimate application writing shortcuts for example. This folder should always be inspected to make sure that all the files in there are legitimate"
        ],
        "id": "2aa0a6b4-a865-495b-ab51-c28249537b75",
        "logsource": {
            "category": "file_event",
            "product": "windows"
        },
        "references": [
            "https://github.com/OTRF/detection-hackathon-apt29/issues/12",
            "https://threathunterplaybook.com/evals/apt29/detections/5.B.1_611FCA99-97D0-4873-9E51-1C1BA2DBB40D.html"
        ],
        "tags": [
            "attack.persistence",
            "attack.t1547.001"
        ]
    },
    {
        "group": "Sigma",
        "kind": "individual",
        "document": {
            "kind": "jsonl",
            "data": {
                "@timestamp": "2022-09-27T12:25:35.936Z",
                "agent": {
                    "ephemeral_id": "f3e67ec5-847a-448b-8ec4-bac040a44990",
                    "hostname": "CLIENT1",
                    "id": "cd699c72-0f46-4cd9-938a-43faf035a44a",
                    "name": "CLIENT1",
                    "type": "winlogbeat",
                    "version": "7.10.2"
                },
                "ecs": {
                    "version": "1.5.0"
                },
                "event": {
                    "action": "File created (rule: FileCreate)",
                    "category": [
                        "file"
                    ],
                    "code": 11,
                    "created": "2022-09-27T12:27:28.842Z",
                    "kind": "event",
                    "module": "sysmon",
                    "provider": "Microsoft-Windows-Sysmon",
                    "type": [
                        "creation"
                    ]
                },
                "file": {
                    "directory": "C:\\Users\\client1\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs",
                    "name": "Startup",
                    "path": "C:\\Users\\client1\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup"
                },
                "host": {
                    "name": "CLIENT1.breach.local"
                },
                "log": {
                    "level": "information"
                },
                "message": "File created:\nRuleName: T1023\nUtcTime: 2022-09-27 12:25:35.936\nProcessGuid: {0e37ca90-ebbc-6332-4600-000000000500}\nProcessId: 3760\nImage: C:\\Windows\\Explorer.EXE\nTargetFilename: C:\\Users\\client1\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\nCreationUtcTime: 2022-09-27 12:25:35.936\nUser: BREACH\\client1",
                "process": {
                    "entity_id": "{0e37ca90-ebbc-6332-4600-000000000500}",
                    "executable": "C:\\Windows\\Explorer.EXE",
                    "name": "Explorer.EXE",
                    "pid": 3760
                },
                "rule": {
                    "name": "T1023"
                },
                "winlog": {
                    "api": "wineventlog",
                    "channel": "Microsoft-Windows-Sysmon/Operational",
                    "computer_name": "CLIENT1.breach.local",
                    "event_data": {
                        "CreationUtcTime": "2022-09-27 12:25:35.936",
                        "User": "BREACH\\client1"
                    },
                    "event_id": 11,
                    "opcode": "Info",
                    "process": {
                        "pid": 2076,
                        "thread": {
                            "id": 2556
                        }
                    },
                    "provider_guid": "{5770385f-c22a-43e0-bf4c-06f5698ffbd9}",
                    "provider_name": "Microsoft-Windows-Sysmon",
                    "record_id": 1365,
                    "task": "File created (rule: FileCreate)",
                    "user": {
                        "domain": "NT AUTHORITY",
                        "identifier": "S-1-5-18",
                        "name": "SYSTEM",
                        "type": "User"
                    },
                    "version": 2
                }
            }
        },
        "name": "Startup Folder File Write",
        "timestamp": "2022-09-27T12:25:35.936+00:00",
        "authors": [
            "Roberto Rodriguez (Cyb3rWard0g)",
            "OTR (Open Threat Research)"
        ],
        "level": "medium",
        "source": "sigma",
        "status": "experimental",
        "falsepositives": [
            "An FP could be caused by legitimate application writing shortcuts for example. This folder should always be inspected to make sure that all the files in there are legitimate"
        ],
        "id": "2aa0a6b4-a865-495b-ab51-c28249537b75",
        "logsource": {
            "category": "file_event",
            "product": "windows"
        },
        "references": [
            "https://github.com/OTRF/detection-hackathon-apt29/issues/12",
            "https://threathunterplaybook.com/evals/apt29/detections/5.B.1_611FCA99-97D0-4873-9E51-1C1BA2DBB40D.html"
        ],
        "tags": [
            "attack.persistence",
            "attack.t1547.001"
        ]
    },
    {
        "group": "Sigma",
        "kind": "individual",
        "document": {
            "kind": "jsonl",
            "data": {
                "@timestamp": "2022-09-27T12:25:35.936Z",
                "agent": {
                    "ephemeral_id": "f3e67ec5-847a-448b-8ec4-bac040a44990",
                    "hostname": "CLIENT1",
                    "id": "cd699c72-0f46-4cd9-938a-43faf035a44a",
                    "name": "CLIENT1",
                    "type": "winlogbeat",
                    "version": "7.10.2"
                },
                "ecs": {
                    "version": "1.5.0"
                },
                "event": {
                    "action": "File created (rule: FileCreate)",
                    "category": [
                        "file"
                    ],
                    "code": 11,
                    "created": "2022-09-27T12:27:28.842Z",
                    "kind": "event",
                    "module": "sysmon",
                    "provider": "Microsoft-Windows-Sysmon",
                    "type": [
                        "creation"
                    ]
                },
                "file": {
                    "directory": "C:\\Users\\client1\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup",
                    "extension": "ini",
                    "name": "desktop.ini",
                    "path": "C:\\Users\\client1\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\desktop.ini"
                },
                "host": {
                    "name": "CLIENT1.breach.local"
                },
                "log": {
                    "level": "information"
                },
                "message": "File created:\nRuleName: T1023\nUtcTime: 2022-09-27 12:25:35.936\nProcessGuid: {0e37ca90-ebbc-6332-4600-000000000500}\nProcessId: 3760\nImage: C:\\Windows\\Explorer.EXE\nTargetFilename: C:\\Users\\client1\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\desktop.ini\nCreationUtcTime: 2022-09-27 12:25:35.936\nUser: BREACH\\client1",
                "process": {
                    "entity_id": "{0e37ca90-ebbc-6332-4600-000000000500}",
                    "executable": "C:\\Windows\\Explorer.EXE",
                    "name": "Explorer.EXE",
                    "pid": 3760
                },
                "rule": {
                    "name": "T1023"
                },
                "winlog": {
                    "api": "wineventlog",
                    "channel": "Microsoft-Windows-Sysmon/Operational",
                    "computer_name": "CLIENT1.breach.local",
                    "event_data": {
                        "CreationUtcTime": "2022-09-27 12:25:35.936",
                        "User": "BREACH\\client1"
                    },
                    "event_id": 11,
                    "opcode": "Info",
                    "process": {
                        "pid": 2076,
                        "thread": {
                            "id": 2556
                        }
                    },
                    "provider_guid": "{5770385f-c22a-43e0-bf4c-06f5698ffbd9}",
                    "provider_name": "Microsoft-Windows-Sysmon",
                    "record_id": 1366,
                    "task": "File created (rule: FileCreate)",
                    "user": {
                        "domain": "NT AUTHORITY",
                        "identifier": "S-1-5-18",
                        "name": "SYSTEM",
                        "type": "User"
                    },
                    "version": 2
                }
            }
        },
        "name": "Startup Folder File Write",
        "timestamp": "2022-09-27T12:25:35.936+00:00",
        "authors": [
            "Roberto Rodriguez (Cyb3rWard0g)",
            "OTR (Open Threat Research)"
        ],
        "level": "medium",
        "source": "sigma",
        "status": "experimental",
        "falsepositives": [
            "An FP could be caused by legitimate application writing shortcuts for example. This folder should always be inspected to make sure that all the files in there are legitimate"
        ],
        "id": "2aa0a6b4-a865-495b-ab51-c28249537b75",
        "logsource": {
            "category": "file_event",
            "product": "windows"
        },
        "references": [
            "https://github.com/OTRF/detection-hackathon-apt29/issues/12",
            "https://threathunterplaybook.com/evals/apt29/detections/5.B.1_611FCA99-97D0-4873-9E51-1C1BA2DBB40D.html"
        ],
        "tags": [
            "attack.persistence",
            "attack.t1547.001"
        ]
    },
    {
        "total_hits": 4
    }
]